# myservice/server/myserver.service.in
[Unit]
Description=MyDaemon

[Service]
Type=simple
ExecStart=@prefix@/bin/mydaemond

[Install]
WantedBy=multi-user.target

